# IWF/package/scripts/incrond.conf
description "incrond"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 10 5

script
exec /usr/local/sbin/incrond -n
end script




