# IWF/package/scripts/snmpd.conf
description "snmpd"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 10 5

script
exec /usr/sbin/snmpd -f
end script



