# eAccelero/scripts/utils/package/memcached.conf
description "Memcached"

start on runlevel [2345]
stop on runlevel [016]

respawn
respawn limit 10 5

script
exec /usr/bin/memcached -p 50000 -m 1024 -u root
end script
