# eAccelero/scripts/utils/package/pre-iwf.conf
description "script which runs before dre script during boot"

start on runlevel [2345]

task
exec echo "no" >/etc/.iwf/boot

task
exec echo "1" >/etc/.iwf/couchbase_state

